# etc/xdg/autostart/fix-alc245-headphone.desktop
[Desktop Entry]
Type=Application
Name=Fix ALC245 Headphone
Exec=/usr/local/bin/fix-alc245-headphone
X-GNOME-Autostart-enabled=true
